# debian/pkgtemplates/flavours.conffiles
/etc/init.d/@FLAVOUR@
/etc/logrotate.d/@FLAVOUR@
/etc/cron.d/@FLAVOUR@
