# debian/pkgtemplates/flavours.cron.d
30 10 * * * root /usr/sbin/logrotate --force /etc/@FLAVOUR@/logrotate.conf
